# crontab
15 * * * * exec /root/deploy/dbox_whitelist.sh

